# policyagent/src/main/application/bin/policyagent-init
#!/bin/bash

### BEGIN INIT INFO
# Provides:          policyagent-init
# Required-Start:    $local_fs $network $remote_fs
# Required-Stop:     $local_fs $network $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: policyagent-init
# Description:       policyagent-init
### END INIT INFO

# Ensure the script is loaded just after init scripts

logfile=/var/log/policyagent-init.log
whoami >> $logfile
tcsd 2>> $logfile
ps -ef | grep tcsd >> $logfile


case "$1" in
 start)
	echo "policyagent-init : Starting up policyagent-init"
	echo "Start call received" >> $logfile
	echo "==========================" >> $logfile
	echo "Initial mounts ===========" >> $logfile
        echo "******************************************" >> $logfile
	mount >> $logfile
        if [ -d "$DISK_LOCATION" ] && [ -d "$ENC_KEY_LOCATION" ]; then
            /usr/bin/python /opt/policyagent/bin/mount_on_reboot.py
        else
            echo "Cannot re mount the directories, enc_disk directories and encrypted keys does not exist" >> $logfile
            exit 1
        fi
        echo "Final mounts ===========" >> $logfile
        echo "******************************************" >> $logfile
        mount >> $logfile
   ;;
 stop)
	echo "Stop call received" >> $logfile
   ;;
 restart)
	echo "Restart call received" >> $logfile
   ;;
 *)
   echo "Usage: policyagent-init {start|stop|restart}" >&2
   exit 3
   ;;
esac

exit 0
